# systemd/piwifi-flask.service
[Unit]
Description=Pi WiFi Provisioning Flask UI
After=network.target NetworkManager.service
Wants=NetworkManager.service
PartOf=piwifi-manager.service

[Service]
Type=simple
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Configuration via environment file
EnvironmentFile=-/etc/default/piwifi

[Install]
WantedBy=multi-user.target