# systemd/piwifi-manager.service
[Unit]
Description=Pi WiFi Manager (Client/AP fallback)
After=NetworkManager.service
Wants=NetworkManager.service piwifi-flask.service

[Service]
Type=simple
ExecStart=/usr/local/bin/piwifi-manager.sh
Restart=on-failure
RestartSec=2

# Configuration via environment file
EnvironmentFile=-/etc/default/piwifi

[Install]
WantedBy=multi-user.target